# base/config/systemd/user/dotfiles.service
[Unit]
Description=Update dotfiles

[Service]
Type=oneshot
ExecStart=%h/dotfiles/update
