# base/config/systemd/user/dotfiles.timer
[Unit]
Description=Update dotfiles

[Timer]
OnUnitActiveSec=6h
Unit=dotfiles.service

[Install]
WantedBy=timers.target
